# etc/systemd/system/update-hosts.service
[Unit]
Description=update-hosts service

[Service]
ExecStart=/opt/scriptfunction/update-hosts.sh
